# agent-browser-wishlist-skill/examples/launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.wishlist-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/ABSOLUTE/PATH/TO/agent-browser-wishlist-skill/scripts/check_all.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/agent-browser-wishlist-skill</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/wishlist-agent.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/wishlist-agent.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
